# .config/autostart/picom.desktop
[Desktop Entry]
Name=picom
Exec=picom --daemon
Type=Application
